# ai/profiles/macbook-m5-64gb/ollama/com.kehle.ollama.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for the Ollama runtime (macOS / Apple Silicon).

  WHY THIS EXISTS (and not `brew services ollama`):
  The Homebrew `ollama` bottle (0.30.7/0.30.8 era) ships ONLY the Go `ollama`
  binary plus an MLX/Metal stub — it does NOT bundle `llama-server`. ollama
  routes GGUF models (our entire ~/.ollama store) through `llama-server`, so the
  brew install cannot load any model and falls back to a CPU path that also
  fails to start a server. The macOS App bundled a working `llama-server`, which
  is why the App used to be the de-facto runtime.

  FIX: install the official standalone runtime tarball
  (https://github.com/ollama/ollama/releases -> ollama-darwin.tgz), which
  contains `ollama` + the full `lib/ollama/` backend incl. `llama-server` and
  the Metal libs, into a user-owned path (no sudo, no GUI). This agent runs that
  binary's `serve` with all tuning env vars baked in.

  Install layout (user-owned, no sudo):
      ~/.local/ollama/bin/ollama          (the binary)
      ~/.local/ollama/lib/ollama/         (llama-server + ggml/metal backend)
      ~/.local/bin/ollama -> ../ollama/bin/ollama   (on PATH)

  Install path for this agent:
      ~/Library/LaunchAgents/com.kehle.ollama.plist
  Manual (re)load:
      cp com.kehle.ollama.plist ~/Library/LaunchAgents/
      launchctl bootout   gui/$(id -u)/com.kehle.ollama 2>/dev/null || true
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.kehle.ollama.plist
      launchctl enable    gui/$(id -u)/com.kehle.ollama

  NOTE: launchd does not expand ~ or $HOME inside <string> values. Replace
  /Users/isaac everywhere below if your username differs. This supersedes the
  old env-only agent com.kehle.ollama-env (env vars are baked in here now).
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.kehle.ollama</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/isaac/.local/ollama/bin/ollama</string>
		<string>serve</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>OLLAMA_HOST</key>
		<string>0.0.0.0:11434</string>
		<key>OLLAMA_MAX_LOADED_MODELS</key>
		<string>3</string>
		<key>OLLAMA_NUM_PARALLEL</key>
		<string>1</string>
		<key>OLLAMA_KEEP_ALIVE</key>
		<string>30m</string>
		<key>OLLAMA_FLASH_ATTENTION</key>
		<string>1</string>
		<key>OLLAMA_KV_CACHE_TYPE</key>
		<string>q8_0</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/isaac/.local/ollama/ollama.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/isaac/.local/ollama/ollama.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/isaac/.local/ollama</string>
</dict>
</plist>
