# _config/profile.d/_1password
# ---- 1Password / op: zero prompts on shell startup ----

# Only do anything if `op` is installed
if command -v op >/dev/null 2>&1; then
  # Only touch 1Password if there is already an authenticated session
  if op account list >/dev/null 2>&1; then
    token="$(op read 'op://dev/Service Account - cli-automation/token' 2>/dev/null)"
    if [[ -n "$token" ]]; then
      export OP_SERVICE_ACCOUNT_TOKEN="$token"
    fi
  fi
fi

# Optional helper: run this *manually* when you actually want to sign in
op_signin() {
  op signin --account my.1password.com
}
