# _config/profile.d/_fabric
# Fabric alias (brew installs as fabric-ai)
command -v fabric-ai >/dev/null 2>&1 && alias fabric='fabric-ai'
