# _config/profile.d/_windsurf
# Windsurf → Devin Desktop compat shim (rebranded June 2026)
# shellcheck source=_devin
. "$(dirname "${BASH_SOURCE[0]}")/_devin"
