# _config/profile.d/_windsurf
# Windsurf
[ -d "$HOME/.codeium/windsurf/bin" ] || return 0
export PATH="$HOME/.codeium/windsurf/bin:$PATH"
