# profile.d/_1password
1|# ---- 1Password / op: suppress permission prompts during startup scripts ----
2|export OP_REQUIRE_UNRESTRICTED_ACCESS=0
3|
4|# Only do anything if `op` is installed
5|if command -v op >/dev/null 2>&1; then
6|  # Only touch 1Password if there is already an authenticated session
7|  if op account list >/dev/null 2>&1; then
8|    token="$(op read 'op://dev/Service Account - cli-automation/token' 2>/dev/null)"
9|    if [[ -n "$token" ]]; then
10|      export OP_SERVICE_ACCOUNT_TOKEN="$token"
11|    fi
12|  fi
13|fi
14|
15|# Optional helper: run this *manually* when you actually want to sign in
16|op_signin() {
17|  op signin --account my.1password.com
18|}
19|