# profile.d/_fabric
1|# Fabric alias (brew installs as fabric-ai)
2|command -v fabric-ai >/dev/null 2>&1 && alias fabric='fabric-ai'
3|