# profile.d/_git
1|export SETTINGS_BASE="$HOME/code/isaackehle/settings" # the github repo's name
2|
3|