# profile.d/_home_assistant
1|
2|export HOMEASSISTANT_URL=YOUR_HOMEASSISTANT_URL
3|export HOMEASSISTANT_TOKEN="YOUR_LONG_LIVED_TOKEN"
4|
5|