# ai/profiles/macbook-intel-2019-16gb/com.kehle.ollama.intel-2019.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.kehle.ollama.intel-2019.plist
  LaunchAgent override for the 2019 Intel MacBook Pro 16" in the garage.

  Differences from the M5 Max agent:
    1. OLLAMA_NUM_GPU=0 — CPU-only. The AMD Radeon Pro 4-8GB dGPU is too
       small for any useful LLM context, and Metal/ROCm support on Intel
       macs is partial. CPU is faster and more reliable.
    2. OLLAMA_NUM_PARALLEL=1 — single-request at a time; 16GB RAM total.
    3. OLLAMA_MAX_LOADED_MODELS=1 — only one resident model; ~10GB usable.
    4. OLLAMA_KV_CACHE_TYPE=q8_0 — same as the M-series; fits in CPU RAM.
    5. OLLAMA_FLASH_ATTENTION=1 — works on CPU, saves KV memory.
    6. OLLAMA_KEEP_ALIVE=30m — same idle eviction as the main box.
    7. OLLAMA_HOST=0.0.0.0:11434 — bind to all interfaces so the main M5
       Max can reach this box over Tailscale (100.x.x.x:11434).

  Install on the garage machine:
    cp com.kehle.ollama.intel-2019.plist ~/Library/LaunchAgents/com.kehle.ollama.plist
    launchctl unload ~/Library/LaunchAgents/com.kehle.ollama.plist 2>/dev/null
    launchctl load -w ~/Library/LaunchAgents/com.kehle.ollama.plist

  Then on the main M5 Max, add a Tailscale-routed model alias in
  ~/.ollama/config.json or via a Tailscale Funnel if you want to share
  inference publicly.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kehle.ollama</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Ollama.app/Contents/Resources/ollama</string>
    <string>serve</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OLLAMA_HOST</key>
    <string>0.0.0.0:11434</string>

    <!-- Intel 2019 specific: CPU-only -->
    <key>OLLAMA_NUM_GPU</key>
    <string>0</string>

    <key>OLLAMA_NUM_PARALLEL</key>
    <string>1</string>

    <key>OLLAMA_MAX_LOADED_MODELS</key>
    <string>1</string>

    <key>OLLAMA_KV_CACHE_TYPE</key>
    <string>q8_0</string>

    <key>OLLAMA_FLASH_ATTENTION</key>
    <string>1</string>

    <key>OLLAMA_KEEP_ALIVE</key>
    <string>30m</string>

    <!-- 16GB RAM → conservative context. qwen3:14b at 32k ctx fits. -->
    <key>OLLAMA_CONTEXT_LENGTH</key>
    <string>32768</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/com.kehle.ollama.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.kehle.ollama.err</string>

  <key>WorkingDirectory</key>
  <string>/tmp</string>
</dict>
</plist>
