# ai/config/profile.d/_devin
# Devin Desktop (formerly Windsurf)
[ -d "$HOME/.codeium/devin/bin" ] && export PATH="$HOME/.codeium/devin/bin:$PATH"
[ -d "$HOME/.codeium/windsurf/bin" ] && export PATH="$HOME/.codeium/windsurf/bin:$PATH"
