# config/profile.d/_ai
# AI workspace paths
export JOB_SEARCH_WORKSPACE="$PRIMARY_VAULT/job_search"
export PRIVATE_CONFIG_ROOT="$PRIMARY_VAULT/private"
