# config/profile.d/_claude_code
# Claude Code — route via LiteLLM or Ollama instead of Anthropic API

# Check if LiteLLM is running (port 4000) — preferred
if curl -s --max-time 2 http://localhost:4000/v1/models >/dev/null 2>&1; then
    export ANTHROPIC_API_KEY="sk-ant-api03-xxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxx"
    export ANTHROPIC_BASE_URL="http://localhost:4000"
    # Optional: default model via LiteLLM
    # export CLAUDE_MODEL="claude-sonnet-4-6"
# Check if Ollama is running (port 11434) — fallback
elif curl -s --max-time 2 http://localhost:11434/api/tags >/dev/null 2>&1; then
    export ANTHROPIC_API_KEY="sk-ant-api03-xxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxx"
    export ANTHROPIC_BASE_URL="http://localhost:11434"
fi
