# config/profile.d/_git
export SETTINGS_BASE="$HOME/code/isaackehle/settings" # the github repo's name

