# config/profile.d/_git
export SETTINGS_REPO="$HOME/code/isaackehle/settings"

