# config/profile.d/_grok
# Grok CLI — route through LiteLLM proxy
command -v grok >/dev/null 2>&1 || return 0
export GROKCLI_PROVIDER=openai
export GROKCLI_BASE_URL=http://localhost:4000/v1
export GROKCLI_API_KEY=sk-local
