# config/profile.d/_home_assistant

export HOMEASSISTANT_URL=YOUR_HOMEASSISTANT_URL
export HOMEASSISTANT_TOKEN="YOUR_LONG_LIVED_TOKEN"

