# link-bot/com.user.link-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.link-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/isaac/code/isaackehle/settings/link-bot/bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/isaac/code/isaackehle/settings/link-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>BOT_TOKEN</key>
        <string>__FILL_IN__</string>
        <key>ALLOWED_USER_ID</key>
        <string>__FILL_IN__</string>
        <key>OLLAMA_HOST</key>
        <string>http://localhost:11434</string>
        <key>OLLAMA_MODEL</key>
        <string>qwen3:14b</string>
        <key>DESTINATIONS</key>
        <string>{"incoming":"/Users/isaac/code/isaackehle/settings/incoming","links":"/Users/isaac/code/isaackehle/settings/saved-links"}</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/link-bot.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/link-bot.err</string>
</dict>
</plist>
